# systemd-unit/pmbs-snapshot.service
[Unit]
Description=Create btrfs snapshot (every minute) for pmbs

[Service]
ExecStart=/usr/bin/pmbs config snapshot
