# systemd-unit/system/pmbs-clean.service
[Unit]
Description=Clean btrfs snapshots for pmbs

[Service]
ExecStart=/usr/bin/pmbs config clean
